# contrib/lokinet-exit.crontab
# lokinet exit cronjobs

SHELL=/bin/bash

* * * * * /usr/local/bin/lokinet-update-exit-address.sh
* * * * * /usr/local/bin/lokinet-kill-scans.sh