# contrib/lokinet-firewall.crontab
# lokinet firewall cronjobs

SHELL=/bin/bash

0 * * * * /usr/local/bin/lokinet-update-firewall.sh